# files/netkeeper.init
#!/bin/sh /etc/rc.common

START=99
STOP=15

start() {
    (/usr/lib/netkeeper/netkeeper.sh &)
    logger -t netkeeper "listen start"
}

stop() {
    if [ -n "$(pgrep netkeeper.sh)" ]; then
        kill "$(pgrep netkeeper.sh)"
    fi
    if [ -n "$(pgrep pppoe-server)" ]; then
        kill "$(pgrep pppoe-server)"
    fi
}
